# Service/hellolinux.service
[Unit]
Description=Hello Linux Service

[Service]
Type=notify
ExecStart=/usr/sbin/hellolinux/HelloLinux

[Install]
WantedBy=multi-user.target